# etc/profile.d/20-io.github.black-desk.debian-tweak.sh
if [[ :$PATH: != *:/opt/io.github.black-desk.debian-tweak/bin:* ]]; then
	PATH=/opt/io.github.black-desk.debian-tweak/bin:$PATH
fi
export PATH
